# systemd/carcharger-tesla-sync.service
[Unit]
Description=CarCharger Tesla Energy Sync Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/carcharger
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/carcharger/.venv/bin/python3 %h/carcharger/tesla_energy_sync.py --config %h/carcharger/config.yaml
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
